# meta/supervisord-fail2ban.conf
[program:fail2ban]
command=service fail2ban start
autostart = true
autorestart = false
stopasgroup=true
