# 005-sqli/supervisord-mysqld.conf
[program:mysqld]
command=/start-mysqld.sh
numprocs=1
autostart=true
autorestart=true
user=mysql
aenvironment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8,HOME="/var/lib/mysql",USER="mysql"
