# 012-command-inject/supervisord-flask.ini
[program:flask]
command=/opt/app/serve.sh
numprocs=1
autostart=true
autorestart=true
environment=HOME="/opt/app",USER="app"
