# 006-sqli-union/supervisord-apache2.conf
[program:apache2]
command=/start-apache2.sh
numprocs=1
autostart=true
utorestart=true
aenvironment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8,HOME="/var/www",USER="www-data"
