# 020-pickle/supervisord-flask.ini
[program:flask]
command=/bin/sh -c "echo 'NOVI{p1ckle_r1ck}' > /tmp/flag.txt && /opt/app/serve.sh"
numprocs=1
autostart=true
autorestart=true
user=app
environment=HOME="/opt/app",USER="app"
