# 018-eval/supervisord-npm.conf
[program:node]
command=/start-npm.sh
numprocs=1
autostart=true
autorestart=true
user=node
environment=HOME="/opt/app",USER="node"
